# conky/conky.service
# ~/.config/systemd/user/conky.service

[Unit]
Description=Conky

[Service]
Type=simple
ExecStart=/usr/bin/conky
Environment=DISPLAY=:0
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
